# nmserver/rootfs/etc/services.d/syslog-ng/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start syslog-ng service if enabled
# ==============================================================================

## If SSH is disabled, use a fake sleep process
#if ! bashio::var.has_value "$(bashio::addon.port 22)"; then
#    bashio::log.warning "SSH port is disabled. Prevent start of SSH server." 
#    exec sleep infinity
#fi

CONFIG_DIR=/config/nmserver

if ! bashio::fs.directory_exists "${CONFIG_DIR}"; then
  bashio::log.info "Copiando configuracao padrao para ${CONFIG_DIR}"
  cp -r /addons/nmserver/rootfs/etc/nginx/ha-config-ori/nmserver /config
fi

bashio::log.info "Starting the syslog-ng daemon..."
exec /usr/sbin/syslog-ng -F -e
